# etc/crontab
* * * * * . "$HOME/.profile"; eval $(display-env); email -S >/dev/null 2>&1
