# etc/crontab
* * * * * . "$HOME/.profile"; eval $(display-env); email -S
