# etc/crontab
* * * * * export PATH="$PATH:$HOME/.local/bin"; eval $(display-env); email -S
