# etc/crontab
* * * * * PATH="$PATH:$HOME/.local/bin:/usr/local/bin" email -S
