# .config/crontab
* * * * * /home/jer/.local/bin/check_mail
* * * * * mbsync -Va
