# .config/crontab
* * * * * /home/jer/.local/bin/check_mail
*/3 * * * * mbsync -Va
