# .config/crontab
* * * * * PATH="$PATH:$HOME/.local/bin" email check cron
