# .config/crontab
* * * * * PATH="$PATH:$HOME/.local/bin:/usr/local/bin" email -Sn
