# .config/crontab
* * * * * PATH="$PATH:$HOME/.local/bin:/usr/local/bin" email check cron 2>&1 >> ~/leedle
